# config/systemd/user/omadora-weather-check.timer
[Unit]
Description=Omadora Weather Check Timer
PartOf=omadora-session.target

[Timer]
OnBootSec=30s
OnUnitActiveSec=5min
AccuracySec=30s
RandomizedDelaySec=15s
