# config/systemd/user/omadora-battery-monitor.service
[Unit]
Description=Omadora Battery Monitor Service

[Service]
Type=oneshot
ExecStart=%h/.local/share/omadora/libexec/systemd/omadora-battery-monitor
