# config/systemd/user/omadora-theme-reload.path
[Unit]
Description=Omadora Theme Reload Path
PartOf=omadora-session.target
After=graphical-session.target

[Path]
PathChanged=%h/.config/omadora/current/theme.name
