# config/systemd/user/omadora-session.target
[Unit]
Description=Omadora Session Target
PartOf=graphical-session.target
After=graphical-session.target xdg-desktop-portal.service

Wants=mako.service
Wants=waybar.service
Wants=hypridle.service
Wants=hyprpaper.service
#Wants=hyprsunset.service

Wants=omadora-session.service
Wants=omadora-polkitagent.service
Wants=omadora-fcitx5.service
Wants=omadora-battery-monitor.timer
Wants=omadora-update-check.timer
Wants=omadora-theme-reload.path
