# config/systemd/user/hyprsunset.service.d/omadora.conf
[Unit]
PartOf=omadora-session.target

[Service]
StandardOutput=null
