# config/systemd/user/waybar.service.d/omadora.conf
[Unit]
PartOf=omadora-session.target
ConditionPathExists=!%h/.local/state/omadora/toggles/waybar-off
